# src/etc/supervisord.d/httpd-bootstrap.conf
[program:httpd-bootstrap]
autorestart = false
autostart = %(ENV_APACHE_AUTOSTART_HTTPD_BOOTSTRAP)s
command = /usr/sbin/httpd-bootstrap
priority = 6
redirect_stderr = true
startsecs = 0
startretries = 0
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
