# src/etc/services-config/supervisor/supervisord.d/httpd-bootstrap.conf
[program:httpd-bootstrap]
priority = 6
command = /usr/sbin/httpd-bootstrap
autostart = %(ENV_APACHE_AUTOSTART_HTTPD_BOOTSTRAP)s
startsecs = 0
startretries = 0
autorestart = false
redirect_stderr = true
stdout_logfile = /var/log/httpd/error_log
stdout_events_enabled = true