# etc/profile.d/enable-rh-php56.sh
#!/usr/bin/env bash

source /usr/bin/scl_source enable rh-php56