# src/etc/profile.d/httpd-startup.sh
#!/usr/bin/env bash

source /usr/sbin/httpd-startup
