# etc/services-config/supervisor/supervisord.d/httpd-bootstrap.conf
[program:httpd-bootstrap]
priority = 6
command = bash -c 'touch /tmp/httpd-bootstrap.lock; /usr/sbin/httpd-bootstrap && rm -f /tmp/httpd-bootstrap.lock'
startsecs = 0
startretries = 0
autorestart = false
redirect_stderr = true
stdout_logfile = /var/log/httpd/error_log
stdout_events_enabled = true
environment = APP_HOME_DIR="%(ENV_APACHE_CONTENT_ROOT)s",APACHE_SUEXEC_USER_GROUP="false",DATE_TIMEZONE="%(ENV_PHP_OPTIONS_DATE_TIMEZONE)s",SERVICE_USER="%(ENV_APACHE_SYSTEM_USER)s",SUEXECUSERGROUP="false",SERVICE_UID="%(ENV_APACHE_HEADER_X_SERVICE_UID)s"